# LICENSE.md/crontab
*/20 * * * * find /home/${USER}/.local/share/okular/docdata/ -mtime +30 -delete
